# products/rappterverse/rappterverse/scripts/com.rappterverse.local-platform.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rappterverse.local-platform</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kodyw/Documents/GitHub/Rappter/rappterverse/scripts/local_platform.sh</string>
        <string>--loop</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kodyw/Documents/GitHub/Rappter/rappterverse</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/kodyw/Documents/GitHub/Rappter/rappterverse/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kodyw/Documents/GitHub/Rappter/rappterverse/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
